# sliceworker.service
[Unit]
Description=Slice Worker Start

[Service]
ExecStart=/opt/sliceup/scripts/workerstart.sh

[Install]
WantedBy=multi-user.target
